# examples/etc/rc.d/busybeed
#!/bin/sh
#
# $OpenBSD: dhcpd,v 1.2 2011/07/08 04:29:54 ajacoutot Exp $

daemon="/usr/sbin/busybeed"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_cmd $1
